# scripts/auto-maintenance/auto-maintenance.timer
[Unit]
Description=Execute weekly maintenance

[Timer]
OnCalendar=Sun
Persistent=true

[Install]
WantedBy=timers.target
