# scripts/auto-maintenance/auto-maintenance.service
[Unit]
Description=auto-maintenance weekly service

[Service]
Type=oneshot
ExecStart=%h/auto-maintenance.sh
